# hypr.service
Unit]
Description=Hyprland
[Service]
ExecStart=/usr/bin/Hyprland
[Install]
WantedBy=default.target
